# deploy/docker/elektrine-compose.service
[Unit]
Description=Elektrine Docker Compose stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/elektrine/app
ExecStart=/opt/elektrine/app/scripts/deploy/docker_deploy.sh --profile caddy --output /opt/elektrine/app/deploy/docker/generated.docker.yml
ExecStop=/usr/bin/docker compose -f /opt/elektrine/app/deploy/docker/generated.docker.yml down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
